# dot_config/systemd/user/huion-kd100d.service
[Unit]
Description=Huion KD100 keydial daemon (raw HID -> uinput)
# Keep retrying forever — the device may be absent at login or unplugged.
StartLimitIntervalSec=0

[Service]
ExecStart=%h/go/bin/huion-kd100d
Restart=always
RestartSec=3s

[Install]
WantedBy=default.target
