# dot_config/systemd/user/sync-shared-datafiles.service
[Unit]
Description=Pull shared-datafiles from remote on login
After=network-online.target
Wants=network-online.target
Before=autocommit-shared-datafiles.service

[Service]
Type=oneshot
WorkingDirectory=%h/Repos/Misc/shared-datafiles
ExecStart=%h/Repos/Misc/shared-datafiles/.automation/sync.sh
TimeoutStartSec=30

[Install]
WantedBy=default.target
